# deploy/systemd/technician-backup-retention.timer
[Unit]
Description=Daily Tecnicos en RD database dump retention

[Timer]
OnCalendar=*-*-* 10:17:00 UTC
Persistent=true
RandomizedDelaySec=15m
AccuracySec=1m
Unit=technician-backup-retention.service

[Install]
WantedBy=timers.target
